# debian/idle-shutdown.prerm
#!/bin/sh
set -e

case "$1" in
  remove|deconfigure)
    systemctl disable --now idle-shutdown.timer >/dev/null 2>&1 || true
    systemctl disable --now idle-rdp-disconnect-watch.service >/dev/null 2>&1 || true
    systemctl disable --now idle-ssh-disconnect-watch.service >/dev/null 2>&1 || true
    ;;
esac

exit 0
